# etc/supervisor.d/celery.ini
[program:celery]
command=.venv/bin/celery worker -E -A catalog.celery --loglevel=INFO
autostart=true
autorestart=true
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=%(here)s/../var/log/%(program_name)s.log
stdout_logfile=%(here)s/../var/log/%(program_name)s.log
