# etc/supervisor.d/subscribe_teko_queue.ini
[program:subscribe_teko_queue]
command=.venv/bin/flask subscribe_teko_queue
autostart=true
autorestart=true
numprocs=1
priority = 900
redirect_stderr=true
stderr_logfile=%(here)s/../var/log/%(program_name)s.log
stdout_logfile=%(here)s/../var/log/%(program_name)s.log